# scripts/opcua-server.service
[Unit]
Description=Raspberry Pi OPC UA Server
After=network.target

[Service]
Type=simple
User=opcua
Group=opcua
WorkingDirectory=/opt/pi-opcua-server
ExecStart=/opt/pi-opcua-server/venv/bin/python -m backend.opcua_server.server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SupplementaryGroups=gpio i2c spi

[Install]
WantedBy=multi-user.target
